# config/optional/boards/cix-p1/npu/_packages/bsp-cli/etc/systemd/system/cix-npu-firstboot.service
[Unit]
Description=CIX NPU Firstboot Setup
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/cix-npu/.firstboot-done

[Service]
Type=oneshot
ExecStart=/usr/libexec/cix-npu-firstboot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
